# board/virt2real/v2/rootfs-overlay/etc/profile.d/osx_term.sh
#!/bin/sh
# osx error in the mc and nano xterm-256color
export TERM="xterm"
